# examples/remote_syslog.supervisor.conf
# Supervisor conf file for stock Ubuntu 11.04 install
[program:remote_syslog]

command=remote_syslog -D
user=root
group=root
autostart=true
autorestart=true
redirect_stderr=true
